# systemd/played.service
[Unit]
Description=played — auto-archive music you actually listen to
Documentation=https://github.com/mellen9999/played
After=dbus.socket
Wants=dbus.socket

[Service]
Type=simple
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
# resolve via PATH so the unit works for both system-wide and ~/.local installs
ExecStart=/usr/bin/env played watch
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=20

# resource limits — tuned for low-RAM systems (e.g. 8GB SBCs)
MemoryHigh=2G
MemoryMax=4G
CPUQuota=50%

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# %h/Music must be writable; if MUSIC_DIR points elsewhere, add a drop-in
# at ~/.config/systemd/user/played.service.d/music-dir.conf with:
#   [Service]
#   ReadWritePaths=/your/custom/path
ReadWritePaths=%h/.local/share/played %h/.local/state/played %h/.cache/played %h/Music
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target
